# site-cookbooks/quova/files/default/quova-server
#!/bin/sh
#
# tomcat@tomcat-ajp-port@  Startup script for an Apache Tomcat Application Server instance
#
# chkconfig: 2345 77 23
#
# description: Startup script for an Apache Tomcat Application Server instance
#
######################################################################
## WARNING: DO NOT EDIT THIS FILE IN PLACE
##   You must edit this file in the svn repository or your changes 
##   will be lost.
######################################################################
#
# Tomcat start and stop script

# Leave blank to run as root (you'll go to hell, though)
USER=quova
GROUP=quova

QUOVA_HOME=/usr/local/quova
QUOVA_START=$QUOVA_HOME/bin/startserver
QUOVA_STOP=$QUOVA_HOME/bin/stopserver
QUOVA_STATUS=$QUOVA_HOME/bin/serverStats

bail()
{
	echo $*
	echo "Usage: $0 {start|stop|restart|status}"
	echo "   You can use multiple commands, for example:"
	echo "   $0 stop start"
	exit 1
}

start()
{
    sh $QUOVA_START
}

stop()
{
    sh $QUOVA_STOP
}

status()
{
    sh $QUOVA_STATUS
}


[ -d "$CATALINA_HOME" ] || bail "Can't find QUOVA_HOME '$QUOVA_HOME'"
[ -n "$1" ] || bail

if [ -n "$USER" ] ; then
	ID=`id -u -n`
	if [ "$ID" = "root" ] ; then
		exec su $USER $0 $*
	elif [ "$ID" != "$USER" ] ; then
		echo "Attempting to sudo to re-run this script"
		exec sudo $0 $*
	fi
fi

while [ -n "$1" ] ; do
	case "$1" in
		start)
			start
			;;
		restart)
			stop
			start
			;;
		stop)
			stop
			;;
                *)
		bail
		;;
	esac
	shift
done

