# site-cookbooks/zenoss/templates/default/zenoss-stack.erb
#!/bin/sh
#
# chkconfig: 2345 80 30 
# description: Zenoss
#
exec /usr/local/zenoss/ctlscript.sh "$@"
