# site-cookbooks/tomcat/templates/default/tomcat.TEMPLATE.erb
#!/bin/sh
#
# tomcat@tomcat-ajp-port@  Startup script for an Apache Tomcat Application Server instance
#
# chkconfig: 2345 77 23
#
# description: Startup script for an Apache Tomcat Application Server instance
#
######################################################################
## WARNING: DO NOT EDIT THIS FILE IN PLACE
##   You must edit this file in the svn repository or your changes 
##   will be lost.
######################################################################
#
# Tomcat start and stop script

# Leave blank to run as root (you'll go to hell, though)
USER=<%= @t_user %>
GROUP=<%= @t_group %>

JAVA_HOME=/usr/local/java
JAVA_OPTS="-server -Xmx2550m -Xms2550m -XX:MaxNewSize=512m -XX:NewSize=512m -XX:SurvivorRatio=3 -XX:TargetSurvivorRatio=80 -XX:MinHeapFreeRatio=10 -XX:PermSize=40m -XX:MaxPermSize=256m -Djava.awt.headless=true"

# If a custom java_opts.conf file exists, override the JAVA_OPTS specified above.
[ -f /var/tomcat/server<%= @ajp_port %>/etc/java_opts.conf ] && . /var/tomcat/server<%= @ajp_port %>/etc/java_opts.conf

# Add config for JMX monitoring
JAVA_OPTS="$JAVA_OPTS -Dzenjmx.rmi.agent.port=<%= @jmx_port %> -javaagent:/usr/local/tomcat/lib/jagent.jar"

# Set cached DNS values to expire after 30 minutes.
# This is necessary when using ELB, which regularly shifts IP address
JAVA_OPTS="$JAVA_OPTS -Dsun.net.inetaddr.ttl=1800"

CATALINA_HOME=/usr/local/tomcat
CATALINA_BASE=/var/tomcat/server<%= @ajp_port %>
CATALINA_OPTS="-Djavax.net.ssl.trustStore=/var/mtm/ssl/truststore.jks -Djavax.net.ssl.trustStorePassword=truststorepassword -Djavax.net.ssl.keyStore=/var/mtm/ssl/client.jks -Djavax.net.ssl.keyStorePassword=clientkeystorepassword"

export JAVA_HOME JAVA_OPTS CATALINA_HOME CATALINA_BASE CATALINA_OPTS

[ -f /etc/lsb-release ] && LSOF=/usr/bin/lsof || LSOF=/usr/sbin/lsof
TOMCAT_SCRIPT="$CATALINA_HOME/bin/catalina.sh"
CATALINA_LOG="$CATALINA_BASE/logs/catalina.out"

bail()
{
	echo $*
	echo "Usage: $0 {start|stop|restart}"
	echo "   You can use multiple commands, for example:"
	echo "   $0 stop start"
	exit 1
}

status()
{
	if [ ! -f "$CATALINA_LOG" ] ; then
		return 1
	fi
	$LSOF $CATALINA_LOG > /dev/null 2>&1
	RET="$?"
	
   	case "$RET" in
	0)
		return 0
		;;
	1)
		return 1
		;;
	*)
		return $RET
		;;
	esac
}

rotatecat()
{
	[ -f "$CATALINA_LOG.5" ] && mv "$CATALINA_LOG.5" "$CATALINA_LOG.6"
	[ -f "$CATALINA_LOG.4" ] && mv "$CATALINA_LOG.4" "$CATALINA_LOG.5"
	[ -f "$CATALINA_LOG.3" ] && mv "$CATALINA_LOG.3" "$CATALINA_LOG.4"
	[ -f "$CATALINA_LOG.2" ] && mv "$CATALINA_LOG.2" "$CATALINA_LOG.3"
	[ -f "$CATALINA_LOG.1" ] && mv "$CATALINA_LOG.1" "$CATALINA_LOG.2"
	[ -f "$CATALINA_LOG" ] && mv "$CATALINA_LOG" "$CATALINA_LOG.1"
}

stoptomcat()
{
	PIDS=
	if [ -f "$CATALINA_LOG" ] ; then
		PIDS=`$LSOF -t $CATALINA_LOG`
	fi
	if [ -n "$PIDS" ] ; then
		echo "Stopping tomcat pids $PIDS"
		kill $PIDS
	else
		echo "No pid found"
	fi
}

killtomcat()
{
	PIDS=
	if [ -f "$CATALINA_LOG" ] ; then
		PIDS=`$LSOF -t $CATALINA_LOG`
	fi
	if [ -n "$PIDS" ] ; then
		echo "Killing tomcat pids $PIDS"
		kill -9 $PIDS
	else
		echo "No pid found"
	fi
}

start()
{
        if [ $1 == "debug" ] ; then
            echo "Starting tomcat in debugger mode"
            JAVA_OPTS="${JAVA_OPTS} -Xdebug -Xrunjdwp:transport=dt_socket,address=8000,server=y,suspend=n"
            export JAVA_OPTS
        else
            echo "Starting tomcat"
        fi
	status
	RET="$?"
	if [ "$RET" = "0" ] ; then
		echo "NOT STARTING: Tomcat is already running"
		return
	fi
	# This is a hack for the Map, causes the attachments directory to be
	# created in a predictable location
	if [ -d "$CATALINA_BASE/webapps/mom" ] ; then
		cd "$CATALINA_BASE/webapps/mom"
	else
		cd "$CATALINA_BASE/temp/"
	fi
	if [ -f "$CATALINA_LOG" ] ; then
		rotatecat && sleep 1
	fi
	if [ -f "$CATALINA_BASE/conf/jaas.config" ] ; then
		CATALINA_OPTS="$CATALINA_OPTS -Djava.security.auth.login.config=$CATALINA_BASE/conf/jaas.config"
		export CATALINA_OPTS
	fi
	echo ${TOMCAT_SCRIPT} start
	${TOMCAT_SCRIPT} start
}

stop()
{
	echo "Stopping tomcat"
	status
	RET="$?"
	if [ "$RET" = "1" ] ; then
		echo "NOT STOPPING: Tomcat is not running"
		return
	fi
	OLD_OPTS="$JAVA_OPTS"
	JAVA_OPTS="-Xmx256M -Xms256M"
	export JAVA_OPTS

	# ${TOMCAT_SCRIPT} stop
        stoptomcat
	for n in 10 9 8 7 6 5 4 3 2 1 ; do
		status
		RET="$?"
		if [ "$RET" != "1" ] ; then
			echo "Tomcat not stopped, waiting ... [$n]"
			sleep 3
		fi
	done
	status
	[ "$RET" = "0" ] && killtomcat

	JAVA_OPTS="$OLD_OPTS"
	export JAVA_OPTS
}

[ -x "$TOMCAT_SCRIPT" ] || bail "Can't find executable '$TOMCAT_SCRIPT'"
[ -d "$CATALINA_HOME" ] || bail "Can't find CATALINA_HOME '$CATALINA_HOME'"
[ -d "$CATALINA_BASE" ] || bail "Can't find CATALINA_BASE '$CATALINA_BASE'"
[ -n "$1" ] || bail

if [ -n "$USER" ] ; then
	ID=`id -u -n`
	if [ "$ID" = "root" ] ; then
		exec su tomcat $0 $*
	elif [ "$ID" != "$USER" ] ; then
		echo "Attempting to sudo to re-run this script"
		exec sudo $0 $*
	fi
fi

while [ -n "$1" ] ; do
	case "$1" in
		start)
			start "normal"
			;;
		start-debug)
			start "debug"
			;;
		restart)
			stop
			start "normal"
			;;
		stop)
			stop
			;;
	kill)
		killtomcat
		;;
	status)
		status
		RET="$?"
		case "$RET" in
		0)
			echo "Tomcat is running"
			;;
		1)
			echo "Tomcat is stopped"
			;;
		*)
			echo "Tomcat state is unknown ($RET)"
			;;
		esac
		;;
	*)
		bail
		;;
	esac
	shift
done

