# site-cookbooks/mms/templates/default/etc/profile.d/mms-cstools.sh.erb
# Environment variables for cs-tools
#
# NOTE: This file is generated by Chef, don't edit it in place or your
# changes will be lost!

# set PATH to include cs-tools and queue-manager
[ -d <%= @cstools_path %> ] && PATH="<%= @cstools_path %>:$PATH"

export PATH
